# services/ozon_bot.service
[Unit]
Description=Ozon telegram bot
After=multi-user.target

[Service]
Type=simple
ExecStart=/home/lnv/soft/ozon_prices/venv/bin/python3 /home/lnv/soft/ozon_prices/bot_run.py
WorkingDirectory=/home/lnv/soft/ozon_prices
Restart=always

[Install]
WantedBy=multi-user.target